# launchd/io.josephbokan.deadlock-watcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.josephbokan.deadlock-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/josephbokan/deadlock-patch-tool-prod/.venv/bin/python3</string>
        <string>watcher.py</string>
        <string>--interval</string>
        <string>900</string>
        <string>--llm</string>
        <string>claude</string>
        <string>--output-dir</string>
        <string>./site/deadlock/updates</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/josephbokan/deadlock-patch-tool-prod</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/josephbokan/deadlock-patch-tool-prod/logs/watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/josephbokan/deadlock-patch-tool-prod/logs/watcher.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
